# cmd/guardmycopy/guardmycopy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
Template LaunchAgent for guardmycopy.

Manual setup (not automatic):
1) Copy this file to ~/Library/LaunchAgents/com.guardmycopy.agent.plist
2) Replace __GUARDMYCOPY_BIN__ with the absolute path to the guardmycopy binary.
3) Replace __LOG_DIR__ with a writable directory for logs.
4) Load it: launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.guardmycopy.agent.plist
5) Verify: launchctl print gui/$(id -u)/com.guardmycopy.agent

Disable/remove:
- launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.guardmycopy.agent.plist
- rm ~/Library/LaunchAgents/com.guardmycopy.agent.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.guardmycopy.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>__GUARDMYCOPY_BIN__</string>
    <string>run</string>
    <string>--interval</string>
    <string>500</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/guardmycopy.out.log</string>

  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/guardmycopy.err.log</string>
</dict>
</plist>
